# scripts/stream-controller.service
[Unit]
Description=YouTube Stream Controller
After=network.target

[Service]
Type=simple
User=stream
WorkingDirectory=/opt/stream-controller
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/etc/stream-controller/.env
ExecStart=/opt/stream-controller/venv/bin/python -m uvicorn controller.main:app --host 0.0.0.0 --port ${STREAM_CONTROLLER_PORT:-8000}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
